# scripts/install/usr/lib/systemd/system/mero-mkfs@.service
[Unit]
Description=Mero mkfs helper for %i service
Before=mero-server@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/mero/mero-mkfs %i
